# agent_deb/DEBIAN/prerm
#!/bin/bash

set -e

systemctl stop agent.service || true
systemctl disable agent.service || true
rm -f /etc/monitoring_agent.conf
